# gcalcli/gcal-reminder.service
[Unit]
Description=Google Calendar Reminder

[Service]
ExecStart=/usr/bin/gcalcli remind --json 11 '/home/fahad/.dotfiles/gcalcli/remind.py %%s'
